# ansible-roles/roles/pg11-primary/templates/daily.cron.j2
0 4 * * * tmpwatch -v 7d /opt/data/pgsql/pg_archive
